# roles/dynamic-dns/templates/dynamic-dns.timer.j2
[Unit]
Description="Timer for the dynamic-dns.service"

[Timer]
Unit=dynamic-dns.service
OnUnitActiveSec={{ dynamic_dns_repeat_interval }}

[Install]
WantedBy=timers.target
