# roles/certbot/templates/certbot.timer.j2
[Unit]
Description="Timer for the {{ certbot_service_name }}.service"

[Timer]
Unit={{ certbot_service_name }}.service
OnUnitActiveSec={{ certbot_repeat_interval }}

[Install]
WantedBy=timers.target
