# roles/frontman/templates/frontman.service.j2
[Unit]
Description=frontman
After=network.service

[Service]
User={{ frontman_user }}
ExecStart=/usr/sbin/nginx -c {{ frontman_nginx_conf_path }}
Restart=always

[Install]
WantedBy=multi-user.target
