# roles/frontman/templates/frontman.service.j2
[Unit]
Description=frontman
After=network.service

[Service]
User=root
ExecStart=nginx -c /home/{{ frontman_user }}/nginx.conf
Restart=always

[Install]
WantedBy=multi-user.target
